# etc/sv/wshd/run
#!/bin/bash
exec 2>&1

exec wshd --run /run
